# debian/tqftpserv.postinst
#!/bin/sh

#DEBHELPER#

rc-update del tqftpserv default || true
rc-update add tqftpserv sysinit || true

exit 0
